# supervisord.conf
[supervisord]
nodaemon=true
# user=upsource

[program:upsource]
command=/opt/upsource/bin/upsource.sh run
